# .config/systemd/user/battery_check.timer
[Unit]
Description=Run battery_check every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=60

[Install]
WantedBy=default.target
