# .config/systemd/user/battery_check.service
[Unit]
Description=Battery Low Notification

[Service]
ExecStart=[Service]
ExecStart=/home/%u/.local/bin/battery_check.sh

