# .config/systemd/user/battery_check.service
[Unit]
Description=Battery Low Notification

[Service]
ExecStart=/home/$USER/.local/bin/battery_check.sh
