# .config/systemd/user/hyprlock-idle.service
[Unit]
Description=Auto lock screen after idle

[Service]
ExecStart=/usr/bin/hyprlock
Restart=always

[Install]
WantedBy=default.target
